# ansible/install_infra/roles/glusterfs_server_install/templates/gluster-exporter_override.service.j2
# {{ ansible_managed }}

[Service]
{% if gluster_exporter_cpu_quota is defined %}
CPUQuota={{ gluster_exporter_cpu_quota }}
{% endif %}
