# ansible/install_infra/roles/backup_restic/templates/systemd/restic_backup_apps.timer.j2
# {{ ansible_managed }}

[Unit]
Description=Backup apps data with restic
Requires=restic_backup_apps.service

[Timer]
Unit=restic_backup_apps.service
OnCalendar={{ restic_backup_apps_timer }}
Persistent=false

[Install]
WantedBy=timers.target
