# ansible/install_infra/roles/backup_restic/templates/systemd/restic_repo_maintenance.timer.j2
# {{ ansible_managed }}

[Unit]
Description=Perform maintenance on restic repositories

[Timer]
Unit=restic_repo_maintenance.service
OnCalendar={{ restic_repo_maintenance_timer }}
Persistent=false

[Install]
WantedBy=timers.target
