# ansible/install_infra/roles/backup_restic/templates/systemd/restic_backup_host.timer.j2
# {{ ansible_managed }}

[Unit]
Description=Backup host data with restic
Requires=restic_backup_host.service

[Timer]
Unit=restic_backup_host.service
OnCalendar={{ restic_backup_hosts_timer }}
Persistent=false

[Install]
WantedBy=timers.target
